# config/includes.chroot_after_packages/usr/lib/systemd/system/cbpp-apt-setup.service
[Unit]
Description=cbpp apt setup service
After=getty@%i.service

[Service]
ExecStart=/usr/lib/cbpp/apt-setup
Type=oneshot

[Install]
WantedBy=multi-user.target
